# deploy/onprem/heimdex-cleanup.timer
[Unit]
Description=Heimdex Shorts Render Cleanup Timer (daily 02:00 UTC)

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
